# udev/99-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/var/lib/lightdm/.Xauthority", RUN+="/opt/monitor_hotplug.sh"

